# etc/init/utama.conf
# service program utama.py supaya jalan ketika boot dan dijalankan kembali secara otomatis ketika crash

start on runlevel [2345]
stop on runlevel [016]
chdir /home/pi/ariagraharegency
exec python utama.py
respawn
